# orac-bonnet-bridge/orac-bonnet-bridge.service
[Unit]
Description=Orac <-> Adafruit OLED Bonnet
After=mec.service

[Service]
ExecStart=/usr/local/bin/OracBonnetBridge.py

[Install]
WantedBy=multi-user.target
